# misc/monitoring/dashboard/conf/supervisord.conf
# Copyright Materialize, Inc. All rights reserved.
#
# Use of this software is governed by the Business Source License
# included in the LICENSE file at the root of this repository.
#
# As of the Change Date specified in that file, in accordance with
# the Business Source License, use of this software will be governed
# by the Apache License, Version 2.0.

[supervisord]
nodaemon = true

[program:prometheus]
command = /bin/prometheus
    --config.file=/etc/prometheus/prometheus.yml
    --storage.tsdb.path=/prometheus
    --storage.tsdb.retention.size=1GB
    --storage.tsdb.retention.time=15d
    --web.console.libraries=/usr/share/prometheus/console_libraries
    --web.console.templates=/usr/share/prometheus/consoles
    --web.enable-admin-api
redirect_stderr = true
stdout_logfile = /run/prefixed/prometheus
stdout_logfile_maxbytes = 0

[program:grafana]
command = grafana-server
    --homepath="%(ENV_GF_PATHS_HOME)s"
    --config="%(ENV_GF_PATHS_CONFIG)s"
    cfg:default.log.mode="console"
    cfg:default.paths.data="%(ENV_GF_PATHS_DATA)s"
    cfg:default.paths.logs="%(ENV_GF_PATHS_LOGS)s"
    cfg:default.paths.plugins="%(ENV_GF_PATHS_PLUGINS)s"
    cfg:default.paths.provisioning="%(ENV_GF_PATHS_PROVISIONING)s"
redirect_stderr = true
stdout_logfile = /run/prefixed/grafana
stdout_logfile_maxbytes = 0

[program:sql_exporter]
command = /bin/sql_exporter -config.file /sql_exporter/sql_exporter.yml
redirect_stderr = true
stdout_logfile = /run/prefixed/sql_prom
stdout_logfile_maxbytes = 0
